# deploy/dev/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/data -c listen_addresses='*'
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql_err.log

[program:minio]
command=/usr/local/bin/minio server /data/minio --console-address ":9001"
environment=MINIO_ROOT_USER="minioadmin",MINIO_ROOT_PASSWORD="minioadmin"
autostart=true
autorestart=true
stdout_logfile=/var/log/minio.log
stderr_logfile=/var/log/minio_err.log

[program:swamp-server]
command=/app/swamp-server
environment=DATABASE_URL="postgres://swamp:swamp@localhost:5432/swamp?sslmode=disable",S3_ENDPOINT="http://localhost:9000",S3_BUCKET="swamp-artifacts",S3_ACCESS_KEY="minioadmin",S3_SECRET_KEY="minioadmin",S3_USE_PATH_STYLE="true",APP_ENV="development",APP_PORT="8080"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/swamp-server.log
stderr_logfile=/var/log/swamp-server_err.log
